# config/s6-overlay/services.d/nginx/run
#!/command/with-contenv sh

exec nginx -g 'daemon off;'
