# config/s6-overlay/services.d/php-fpm/run
#!/command/with-contenv sh

exec php-fpm -F
