# config/s6-overlay/services.d/crond/run
#!/command/with-contenv sh

exec s6-setuidgid www-data /usr/sbin/crond -f -l 2
