# config/s6-overlay/services.d/crond/run
#!/command/with-contenv sh

exec /usr/sbin/crond -f -l 2
